# default/systemd/ssh-agent.service.d/override.conf
# Configure ssh-agent for systemd socket activation
[Service]
# Set the socket path to match ssh-agent.socket
Environment="SSH_AUTH_SOCK=%t/ssh-agent.socket"
